# crontab
# back up dotfiles and private dotfiles every 10 minutes
*/10 * * * * sh /home/thor/cron/files.sh 2>> /home/thor/cron/logs/files
*/10 * * * * sh /home/thor/cron/private.sh 2>> /home/thor/cron/logs/private
*/10 * * * * sh /home/thor/cron/setup.sh 2>> /home/thor/cron/logs/setup

