# crontab
# note: verify user $HOME is /home/thor
# every hour
0 * * * * /usr/bin/fish /home/thor/.cron/scripts/clean_home.fish 

# every day at 2pm
0 14 * * * /usr/bin/fish /home/thor/.cron/scripts/git_sync_directories.fish 
0 14 * * * /usr/bin/fish /home/thor/.cron/scripts/flatpak_update.fish
0 14 * * * /usr/bin/fish /home/thor/.cron/scripts/rustup_update.fish
# only boorb:
* 14 * * * /usr/bin/fish /home/thor/.cron/scripts/github_name.fish 

